# scripts/systemd/etc/systemd/system/halcs-be-ioc@.service.in
[Unit]
Description=BPM BE IOC %p%I
After=rc-local.service
Requires=halcs-be@%i.service
After=halcs-be@%i.service
PartOf=halcs-ioc@%i.target

[Service]
# Source environment
EnvironmentFile=/etc/sysconfig/bpm-epics-ioc
EnvironmentFile=/etc/sysconfig/bpm-epics-ioc-slot-mapping
Environment=BPM_NUMBER=%i
Environment=BPM_TRY_READ=3
# Execute pre with root
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /var/log/procServ/%p%i
ExecStartPre=/bin/mkdir -p /var/run/procServ/%p%i
ExecStartPre=/bin/sh -c " \
    INSTANCE_IDX=$$(echo ${BPM_NUMBER} | sed 's|.*-||g'); \
    BOARD_IDX=$$(expr ${INSTANCE_IDX} / 2 + $${INSTANCE_IDX} % 2); \
    HALCS_IDX=$$(expr 1 - ${INSTANCE_IDX} % 2); \
    /opt/epics/startup/ioc/bpm-epics-ioc/iocBoot/iocBPM/BPMCheckInitTries.sh ipc:///tmp/malamute $${BOARD_IDX} $${HALCS_IDX} $${BPM_TRY_READ} \
"
WorkingDirectory=<INSTALL_PREFIX>/<IOC_NAME>/iocBoot/iocBPM
# Run procServ with user ioc
ExecStart=/usr/local/bin/procServ -f -n %p%i -i ^C^D ${PROCSERV_BE_PORT_PREFIX}%i ./runBPM.sh ${BPM_ENDPOINT} ${BPM_NUMBER}

[Install]
WantedBy=halcs-ioc@%i.target
